# jobs.cron
*/5 * * * * (cd /app && . ./.cronenv && poetry run python manage.py sendreminders) > /proc/1/fd/1 2> /proc/1/fd/2
