# jobs.cron
* * * * * echo "TODO https://django-cron.readthedocs.io" > /proc/1/fd/1 2> /proc/1/fd/2
